# cron/crontab

# m h  dom mon dow   command

0 * * * * /usr/bin/newsboat -c /home/alex/.local/share/newsboatcache.db -x reload
*/15 * * * * DISPLAY=:0 /home/alex/bin/wallpaper-random

